# debian/postinst
#!/bin/sh

adduser --system nhome

#update-rc.d nhome defaults
#invoke-rc.d nhome start

exit 0

